# baidu-index-hunter-backend/baidu_index_api.service
[Unit]
Description=Baidu Index Hunter API Service
After=network.target

[Service]
User=YOUR_USERNAME
WorkingDirectory=/path/to/BaiduIndexHunter
Environment="PATH=/path/to/python/bin"
ExecStart=/path/to/python/bin/gunicorn -w 4 -b 0.0.0.0:5000 api_server:app --timeout 120
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 